# nightlight_tui/assets/nightlight-schedule.timer
[Unit]
Description=Check nightlight schedule every minute

[Timer]
OnCalendar=*:*:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
