# nightlight_tui/assets/nightlight-schedule.service
[Unit]
Description=Apply nightlight schedule

[Service]
Type=oneshot
ExecStart=__NIGHTLIGHT_SCRIPT__ schedule
